# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
childlogdir=/var/log/supervisor
strip_ansi=true

[program:crecos_server]
command=gunicorn app:app -b 0.0.0.0:8000 --log-level 'debug' --timeout 120 --workers=3 --threads=3 --worker-connections=10
autostart=true
autorestart=true
stdout_logfile=/tmp/crecos.log
stderr_logfile=/tmp/crecos_error.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[inet_http_server]
port = *:9001
username=crecos
password=datavaders2022